# config/com.familycom.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.familycom.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>familycomd</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/familycomd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/familycomd.err.log</string>
</dict>
</plist>
